# homeassistant/addons/roamcore-traccar-proxy/rootfs/etc/services.d/roamcore-traccar-proxy/run
#!/usr/bin/with-contenv bash
set -euo pipefail

/usr/bin/roamcore-traccar-proxy

