# scripts/backups/systemd/roamcore-rclone-sync.service
[Unit]
Description=RoamCore: rclone sync Proxmox backups to Google Drive
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment=SRC_DIR=/var/lib/vz/dump
Environment=DST=gdrive:RoamCore-Backups/proxmox-vzdump
ExecStart=/usr/local/sbin/roamcore-rclone-sync.sh
