# homeassistant/addons/roamcore-traccar-init/rootfs/etc/services.d/roamcore-traccar-init/run
#!/usr/bin/with-contenv bash
set -e
/usr/bin/roamcore-traccar-init
# one-shot: stop the service container after init
s6-svscanctl -t /var/run/s6/services
